# crontab
* * * * * /bin/bash /var/www/html/uploadFileToS3.sh

